# compose/cryproblem/src/supervisord.conf
[supervisord]
nodaemon=true

[program:app]
directory=/root/app
command=/usr/local/bin/npm start

