# compose/cryproblem/src/supervisord.conf
[supervisord]
nodaemon=true

[program:app1]
directory=/root/app1
command=/usr/local/bin/npm start

[program:app2]
directory=/root/app2
command=/usr/local/bin/npm start
